# udev/20-hotplug-monitor.rules
KERNEL=="card0", SUBSYSTEM=="drm", RUN+="/usr/bin/su - chromium -c /home/chromium/display-hotplug.sh"